# qemu/rc.local
#!/bin/bash
set -ex

# Enable unlimited size for coredumps
ulimit -c unlimited

# 9P filesystem mount (example)
mount -t 9p -o trans=virtio,version=9p2000.L,msize=10485760 host0 /mnt/shared
